# v2_6_1/build/configs/supervisord/orocrm.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stopsignal=TERM
stopasgroup=true
autostart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true


[program:php-fpm]
command=/usr/sbin/php-fpm -F
stopasgroup=true
stopsignal=TERM
autostart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true


[program:cron]
command=/usr/sbin/cron -f
stopsignal=TERM
stopasgroup=true
autostart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:oro_web_socket]
directory=/home/orocrm/app/web
command=php ./app/console clank:server --env=prod
numprocs=1
autostart=true
autorestart=true
user=orocommerce
redirect_stderr=true

[program:oro_message_consumer]
directory=/home/orocrm/app/web
command=php ./app/console oro:message-queue:consume --env=prod
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
user=orocommerce
redirect_stderr=true
