# v4_1_4/build/configs/supervisord/orocrm.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stopsignal=TERM
stopasgroup=true
autostart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stderr_logfile_maxbytes=5MB
stdout_logfile=/var/log/supervisor/nginx.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=7


[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
stopasgroup=true
stopsignal=TERM
autostart=true
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stderr_logfile_maxbytes=5MB
stdout_logfile=/var/log/supervisor/php-fmp.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=7


[program:cron]
command=/usr/sbin/cron -f
stopsignal=TERM
stopasgroup=true
autostart=true
stderr_logfile=/var/log/supervisor/cron.err.log
stderr_logfile_maxbytes=5MB
stdout_logfile=/var/log/supervisor/cron.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=7


[program:oro_web_socket]
directory=/var/www/html/orocrm
command=php ./bin/console gos:websocket:server --env=prod
numprocs=1
autostart=true
autorestart=true
user=www-data
redirect_stderr=true

[program:oro_message_consumer]
directory=/var/www/html/orocrm
command=php ./bin/console oro:message-queue:consume --env=prod
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
user=www-data
redirect_stderr=true