# v4_1_4/build/configs/supervisord/supervisord.conf
[supervisord]
nodaemon = true
loglevel = INFO
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)


[program:app]
command = php-fpm
process_name=%(program_name)s-%(process_num)02d
priority = 1
autorestart = unexpected
stopasgroup = true
#stdout_logfile=/dev/fd/1
#stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
#redirect_stderr=true

[program:oro_web_socket]
directory=/var/www/backend
command=php ./bin/console gos:websocket:server --env=prod
numprocs=1
autostart=true
autorestart=true
user=www-data
#redirect_stderr=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:oro_message_consumer]
directory=/var/www/backend
command=php ./bin/console oro:message-queue:consume --env=prod
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
user=www-data
#redirect_stderr=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr